# root/etc/services.d/shiori/run
#!/usr/bin/with-contenv bash

cd /go/bin || exit

exec \
	s6-setuidgid abc /go/bin/shiori \
    serve
